# net-p2p/i2p/files/i2p.service.in
[Unit]
Description=I2P (Invisible Internet Protocol) daemon
After=network.target

[Service]
Environment="PIDDIR=%I2PPIDDIR%" "I2P=%I2PDIR%"
#EnvironmentFile=/etc/conf.d/i2p
PIDFile=%I2PPIDDIR%/i2p.pid
User=i2p
Group=i2p
ExecStart=%I2PDIR%/i2psvc %I2PCONFDIR%/wrapper.config wrapper.daemonize=FALSE wrapper.pidfile=%I2PPIDDIR%/i2p.pid

[Install]
WantedBy=multi-user.target
