# net-p2p/i2p/files/i2p.initd.in
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

PIDDIR=%I2PPIDDIR%
PIDFILE=${PIDDIR}/i2p.pid
I2P=%I2PDIR%
WRAPPER_CMD="$I2P/i2psvc"
WRAPPER_CONF="%I2PCONFDIR%/wrapper.config"

depend() {
	need net
}

start() {
	ebegin "Starting I2P"

	checkpath -q -d -m 0755 -o i2p:i2p "${PIDDIR}"

	start-stop-daemon --start -b --user i2p --group i2p \
		--quiet --wait 500 --exec "$WRAPPER_CMD" "$WRAPPER_CONF" \
		wrapper.daemonize=FALSE wrapper.pidfile="${PIDFILE}"
	eend $?
}

stop() {
	ebegin "Stopping I2P"
	start-stop-daemon --stop --user i2p --group i2p --pidfile "${PIDFILE}"
	eend $?
}
