# app-misc/memlockd/files/memlockd.service
[Unit]
Description=Keep critical files locked in memory

[Service]
Type=forking
ExecStart=/usr/sbin/memlockd

[Install]
WantedBy=multi-user.target
